# crontab
# Remove old nvim undo files
0 12 * * * find $HOME/.local/state/nvim/undo -mtime +28 -exec rm {} \;
