# crontab
# Remove old nvim undo files
0 0 * * 0 find $HOME/.config/nvim/undo -mtime +28 -exec rm {} \;
